# data/gnome-shell.target
[Unit]
Description=GNOME Shell
DefaultDependencies=false

BindsTo=gnome-shell.service
After=gnome-shell.service